# api_server/configs/crontab/zc
# Restart once a day to periodically update manifest.json
0 1 * * * root /bin/systemctl try-restart zc
